# src/permission-agent/harbour-owncloud-permission-agent.service
[Unit]
Description=ownCloud photo backup
After=lipstick.service

[Service]
Type=simple
ExecStart=/usr/bin/harbour-owncloud-permission-agent
Restart=always

[Install]
WantedBy=lipstick.service
